# backup/crontab
00 * * * * aws s3 cp /sqlite-data/data.sqlite s3://csud-reservation-backup/
